# scripts/systemd/genesis-disk-hygiene.timer.template
[Unit]
Description=Genesis Disk Hygiene — daily schedule

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
AccuracySec=1h
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
