# scripts/systemd/genesis-watchdog.timer.template
[Unit]
Description=Genesis Watchdog — periodic health check

[Timer]
OnBootSec=60
OnUnitActiveSec=300
AccuracySec=10

[Install]
WantedBy=timers.target
