# config/genesis-guardian.service
[Unit]
Description=Genesis Guardian health check
Documentation=https://github.com/YOUR_GITHUB_USER/GENesis
# No After/Wants — intentionally standalone. Timer triggers this.

[Service]
Type=oneshot
ExecStart=%h/.local/share/genesis-guardian/.venv/bin/python -m genesis.guardian
WorkingDirectory=%h/.local/share/genesis-guardian
Environment=PYTHONPATH=%h/.local/share/genesis-guardian/src
Environment=GUARDIAN_CONFIG=%h/.local/share/genesis-guardian/config/guardian.yaml
Environment=GUARDIAN_SECRETS=%h/.local/share/genesis-guardian/secrets.env
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin

StandardOutput=journal
StandardError=journal

# CC diagnosis is agentic (multi-turn, tool use) — can take 30-60 min for
# deep investigation + recovery + verification. 65 min allows CC its full
# hour plus 5 min for signal collection, alerts, and cleanup.
TimeoutStartSec=3900

# Resource limits — the Guardian is a LAST-DITCH host-side recovery tool. When
# it runs, Genesis is presumed down and NOTHING revives the Guardian if it dies,
# so it gets GENEROUS resources: the cap is only a host-survival backstop, not a
# workload shaper. MemoryMax=80% of host RAM (systemd resolves the % at runtime)
# is 2-8x what a CC diagnosis needs, while still leaving headroom so a runaway
# Guardian dies at the cap instead of freezing the whole host — which would kill
# the very timer that re-fires it. The real "don't fight a healthy Genesis" guard
# is runtime-aware and lives in diagnosis.py's pre-flight (it refuses CC diagnosis
# when host memory is already tight), not this ceiling.
# MemorySwapMax=0 prevents I/O thrashing under pressure.
# OOMScoreAdjust=0 (neutral) — the Guardian is NOT expendable (it is the last
# line of defense), but must not preferentially kill the user's host-side Claude
# Code or the container either; it competes as an equal under host pressure.
MemoryMax=80%
MemorySwapMax=0
OOMScoreAdjust=0

# I/O limits — prevent Guardian's own diagnosis (CC session, log collection)
# from contributing to I/O saturation on the genesis storage device.
# NOTE: Device path is install-specific. install_guardian.sh detects this
# automatically. The values below are for a typical Incus genesis LV.
# IOReadBandwidthMax=/dev/ubuntu-vg/genesis-lv 50M
# IOWriteBandwidthMax=/dev/ubuntu-vg/genesis-lv 20M
