# config/genesis-guardian.service
[Unit]
Description=Genesis Guardian health check
Documentation=https://github.com/YOUR_GITHUB_USER/GENesis
# No After/Wants — intentionally standalone. Timer triggers this.

[Service]
Type=oneshot
ExecStart=%h/.local/share/genesis-guardian/.venv/bin/python -m genesis.guardian
WorkingDirectory=%h/.local/share/genesis-guardian
Environment=PYTHONPATH=%h/.local/share/genesis-guardian/src
Environment=GUARDIAN_CONFIG=%h/.local/share/genesis-guardian/config/guardian.yaml
Environment=GUARDIAN_SECRETS=%h/.local/share/genesis-guardian/secrets.env
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin

StandardOutput=journal
StandardError=journal

# CC diagnosis is agentic (multi-turn, tool use) — can take 30-60 min for
# deep investigation + recovery + verification. 65 min allows CC its full
# hour plus 5 min for signal collection, alerts, and cleanup.
TimeoutStartSec=3900

# Resource limits — prevent CC diagnostic from consuming unbounded host memory.
# 6G covers the guardian Python process + one CC session + tool subprocesses.
# MemorySwapMax=0 prevents I/O thrashing under pressure.
# OOMScoreAdjust=500 (positive) makes guardian first to be OOM-killed vs
# genesis (-500) or ollama — guardian is expendable, primary workloads are not.
MemoryMax=6G
MemorySwapMax=0
OOMScoreAdjust=500
