# config/genesis-bridge.service
[Unit]
Description=Genesis Telegram Bridge (CC conversation relay)
After=network.target
# Stop crash-looping: max 4 restarts in 120s, then give up until manual reset.
# Must be in [Unit], not [Service] — systemd 255 ignores them in [Service].
StartLimitBurst=4
StartLimitIntervalSec=120

[Service]
Type=simple
Environment=GENESIS_REPO_ROOT=%h/genesis
Environment=VENV_PATH=%h/genesis/.venv
WorkingDirectory=%h
ExecStart=/bin/bash -c 'cd "${GENESIS_REPO_ROOT}" && exec "${VENV_PATH}/bin/python" -m genesis.channels.bridge'
Restart=on-failure
RestartSec=10

# Exit code 200 = "another instance holds the lock" — don't restart for that.
# See genesis.util.process_lock.EXIT_ALREADY_RUNNING.
RestartPreventExitStatus=200

Environment=PATH=%h/genesis/.venv/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
